# jibri/rootfs/etc/services.d/20-icewm0/run
#!/usr/bin/with-contenv bash

DAEMON="/usr/bin/icewm-session --notray --display=:0"
exec s6-setuidgid jibri /bin/bash -c "DISPLAY=:0 exec $DAEMON"

