# jibri/rootfs/etc/services.d/21-icewm1/run
#!/usr/bin/with-contenv bash

DAEMON="/usr/bin/icewm-session --notray --display=:1"
exec s6-setuidgid jibri /bin/bash -c "HOME=/home/jibri DISPLAY=:1 exec $DAEMON"

