# web/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec openresty -c /config/nginx/nginx.conf
